# usb-manager/99-usb-vid-pid-check.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/etc/udev/usb_check.sh %s{idVendor} %s{idProduct} $env{DEVPATH} %k"
